# deploy/camerad.service.in
[Unit]
Description=Panthera-WAM RealSense D405 capture service
After=default.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@
Environment=PYTHONUNBUFFERED=1
Environment=NO_PROXY=127.0.0.1,localhost,::1
Environment=no_proxy=127.0.0.1,localhost,::1
UnsetEnvironment=http_proxy https_proxy HTTP_PROXY HTTPS_PROXY ALL_PROXY all_proxy grpc_proxy GRPC_PROXY
EnvironmentFile=-%h/.config/panthera-wam/armd.env
ExecStart=/usr/bin/env uv run --package panthera-armd camerad --mode auto --bind [::1]:50052
Restart=on-failure
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=5

[Install]
WantedBy=default.target
